# signaling-app-rust/signaling-app-rust.service
[Unit]
Description=Rust signaling server

[Service]
Type=simple
ExecStart=signaling-app-rust

[Install]
WantedBy=multi-user.target