# signaling-app-rust/signaling-app-rust.service
[Unit]
Description=Rust signaling server

[Service]
Type=simple
ExecStart=rust-signaling-server

[Install]
WantedBy=multi-user.target